# trader.conf
[program:trader]
autorestart=true
command=python start_trader.py --pythonpath="/app/bots"
